# BuildAndDeployment/CI/cookbooks/de_ci_web/templates/default/bloq.conf.erb
description "Bloq webapp server"
author "de_ci@ucloud"

start on (filesystem and net-device-up IFACE!=lo)
stop on runlevel [016]

respawn

exec start-stop-daemon --chuid web_bloq --chdir <%= @web_directory %> --start --exec bin/rails -- server -b <%= node['ipaddress'] %> -p <%= node[:paas_agent][:webapp_port] %>

